# arch/etc/udev/rules.d/99-lowbat.rules
# Suspend the system when battery level drop triggers an event
SUBSYSTEM=="power_supply", ATTR{status}=="Discharging", RUN+="/usr/bin/systemctl suspend"
