# desktop/arch/etc/systemd/system/docker.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --insecure-registry 172.30.0.0/16 --data-root=/home/j/docker
